# onboard/rov-arm-boot.service
[Unit]
Description=ROV arm boot guard — hold J6/AUX neutral until new_ar.py starts
After=local-fs.target
Before=network.target

[Service]
Type=simple
WorkingDirectory=/home/uruc/URUCDreadYachet
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -u onboard/arm_boot_guard.py
Restart=on-failure
RestartSec=3
# Release serial quickly when stopped (before MAVProxy from the UI).
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
